# fabfile/templates/django.conf
description "Django instance"
env LANG=en_US.UTF-8
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5
setuid %(user)s
setgid %(user)s
exec %(gunicorn_script)s
